# systemd/user/pw-watchdog.service
[Unit]
Description=Playwright Watchdog Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/detector-rdpbridge/playwright_watchdog.py
Restart=on-failure
RestartSec=5

# Load environment from file (optional)
# EnvironmentFile=%h/.config/pw-watchdog/config.env

# Environment variables
Environment="PW_WATCHDOG_STDOUT=1"
Environment="PW_WATCHDOG_POLL_INTERVAL=0.5"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pw-watchdog

[Install]
WantedBy=default.target


